# onyx_proj/config/services/uat/onyx_config.conf
[unix_http_server]
file=/dev/shm/supervisor.sock   ; (the path to the socket file)

[supervisord]
logfile=/opt/logs/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=1024MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=1           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/dev/shm/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=100                  ; (min. avail startup file descriptors;default 1024)
minprocs=500                 ; (min. avail process descriptors;default 200)
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///dev/shm/supervisor.sock ; use a unix:// URL  for a unix socket

[supervisord]
nodaemon=true

[program:celery_query_executor]
environment="NEW_RELIC_ENVIRONMENT=celery"
environment="NEW_RELIC_CONFIG_FILE=/etc/newrelic/newrelic_onyx.ini"
environment=
         C_FORCE_ROOT=true
command=bash -c "export NEW_RELIC_ENVIRONMENT=celery && export NEW_RELIC_CONFIG_FILE=/etc/newrelic/newrelic_onyx.ini && /usr/local/python3.8/bin/newrelic-admin run-program /usr/local/python3.8/bin/celery worker --time-limit=300 --concurrency=4 --app=celery_app --loglevel=INFO --maxtasksperchild=3000 -Q celery_query_executor"
directory=/usr/local/creditas/Onyx/onyx_proj/onyx_proj/
autostart=true
autorestart=true
startsecs=10                ; number of secs prog must stay running (def. 10)
startretries=3              ; max # of serial start failures (default 3)
exitcodes=1
stopsignal=TERM             ; signal used to kill process (default TERM)
stopwaitsecs=600             ; max num secs to wait before SIGKILL (default 10)
user=root               ; setuid to this UNIX account to run the program
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/1
stderr_logfile_maxbytes=0

[program:celery_callback_resolver]
environment="NEW_RELIC_ENVIRONMENT=celery"
environment="NEW_RELIC_CONFIG_FILE=/etc/newrelic/newrelic_onyx.ini"
environment=
         C_FORCE_ROOT=true
command=bash -c "export NEW_RELIC_ENVIRONMENT=celery && export NEW_RELIC_CONFIG_FILE=/etc/newrelic/newrelic_onyx.ini && /usr/local/python3.8/bin/newrelic-admin run-program /usr/local/python3.8/bin/celery worker --time-limit=300 --concurrency=1 --app=celery_app --loglevel=INFO --maxtasksperchild=3000 -Q celery_callback_resolver"
directory=/usr/local/creditas/Onyx/onyx_proj/onyx_proj/
autostart=true
autorestart=true
startsecs=10                ; number of secs prog must stay running (def. 10)
startretries=3              ; max # of serial start failures (default 3)
exitcodes=1
stopsignal=TERM             ; signal used to kill process (default TERM)
stopwaitsecs=600             ; max num secs to wait before SIGKILL (default 10)
user=root               ; setuid to this UNIX account to run the program
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/1
stderr_logfile_maxbytes=0

[program:campaign_approval]
environment="NEW_RELIC_ENVIRONMENT=celery"
environment="NEW_RELIC_CONFIG_FILE=/etc/newrelic/newrelic_onyx.ini"
environment=
         C_FORCE_ROOT=true
command=bash -c "export NEW_RELIC_ENVIRONMENT=celery && export NEW_RELIC_CONFIG_FILE=/etc/newrelic/newrelic_onyx.ini && /usr/local/python3.8/bin/newrelic-admin run-program /usr/local/python3.8/bin/celery worker --time-limit=300 --concurrency=1 --app=celery_app --loglevel=INFO --maxtasksperchild=3000 -Q celery_campaign_approval"
directory=/usr/local/creditas/Onyx/onyx_proj/onyx_proj/
autostart=true
autorestart=true
startsecs=10                ; number of secs prog must stay running (def. 10)
startretries=3              ; max # of serial start failures (default 3)
exitcodes=1
stopsignal=TERM             ; signal used to kill process (default TERM)
stopwaitsecs=600             ; max num secs to wait before SIGKILL (default 10)
user=root               ; setuid to this UNIX account to run the program
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/1
stderr_logfile_maxbytes=0


[program:celery_click_data]
environment="NEW_RELIC_ENVIRONMENT=celery"
environment="NEW_RELIC_CONFIG_FILE=/etc/newrelic/newrelic_onyx.ini"
environment=
         C_FORCE_ROOT=true
command=bash -c "export NEW_RELIC_ENVIRONMENT=celery && export NEW_RELIC_CONFIG_FILE=/etc/newrelic/newrelic_onyx.ini && /usr/local/python3.8/bin/newrelic-admin run-program /usr/local/python3.8/bin/celery worker --time-limit=300 --concurrency=5 --app=celery_app --loglevel=INFO --maxtasksperchild=3000 -Q celery_click_data"
directory=/usr/local/creditas/Onyx/onyx_proj/onyx_proj/
autostart=true
autorestart=true
startsecs=10                ; number of secs prog must stay running (def. 10)
startretries=3              ; max # of serial start failures (default 3)
exitcodes=1
stopsignal=TERM             ; signal used to kill process (default TERM)
stopwaitsecs=600             ; max num secs to wait before SIGKILL (default 10)
user=root               ; setuid to this UNIX account to run the program
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/1
stderr_logfile_maxbytes=0

[program:celery_campaign_query_executor]
environment="NEW_RELIC_ENVIRONMENT=celery"
environment="NEW_RELIC_CONFIG_FILE=/etc/newrelic/newrelic_onyx.ini"
environment=
         C_FORCE_ROOT=true
command=bash -c "export NEW_RELIC_ENVIRONMENT=celery && export NEW_RELIC_CONFIG_FILE=/etc/newrelic/newrelic_onyx.ini && /usr/local/python3.8/bin/newrelic-admin run-program /usr/local/python3.8/bin/celery worker --time-limit=3000 --concurrency=4 --app=celery_app --loglevel=INFO --maxtasksperchild=900 -n celery_query_executor@%%h -Q celery_campaign_query_executor"
directory=/usr/local/creditas/Onyx/onyx_proj/onyx_proj/
autostart=true
autorestart=true
startsecs=10                ; number of secs prog must stay running (def. 10)
startretries=3              ; max # of serial start failures (default 3)
exitcodes=1
stopsignal=TERM             ; signal used to kill process (default TERM)
stopwaitsecs=600             ; max num secs to wait before SIGKILL (default 10)
user=root               ; setuid to this UNIX account to run the program
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/1
stderr_logfile_maxbytes=0


[program:onyx]
environment="NEW_RELIC_ENVIRONMENT=web"
environment="NEW_RELIC_CONFIG_FILE=/etc/newrelic/newrelic_onyx.ini"
environment=
         C_FORCE_ROOT=true
command=bash -c "export NEW_RELIC_ENVIRONMENT=web && export NEW_RELIC_CONFIG_FILE=/etc/newrelic/newrelic_onyx.ini && /usr/local/python3.8/bin/newrelic-admin run-program /usr/local/python3.8/bin/uwsgi --ini /etc/onyx_uwsgi.ini"
autostart=true
autorestart=true
logfile_maxbytes=1024MB        ; max # logfile bytes b4 rotation (default 50MB)
logfile_backups=3
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/1
stderr_logfile_maxbytes=0
exitcodes=1

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/1
stderr_logfile_maxbytes=0


[program:onyx_celery_callback]
environment="NEW_RELIC_ENVIRONMENT=celery"
environment="NEW_RELIC_CONFIG_FILE=/etc/newrelic/newrelic_onyx.ini"
environment=
         C_FORCE_ROOT=true
command=bash -c "export NEW_RELIC_ENVIRONMENT=celery && export NEW_RELIC_CONFIG_FILE=/etc/newrelic/newrelic_onyx.ini && /usr/local/python3.8/bin/newrelic-admin run-program /usr/local/python3.8/bin/celery worker --time-limit=300 --concurrency=5 --app=celery_app --loglevel=INFO --maxtasksperchild=3000 -Q onyx_celery_callback"
directory=/usr/local/creditas/Onyx/onyx_proj/onyx_proj/
autostart=true
autorestart=true
startsecs=10                ; number of secs prog must stay running (def. 10)
startretries=3              ; max # of serial start failures (default 3)
exitcodes=1
stopsignal=TERM             ; signal used to kill process (default TERM)
stopwaitsecs=600             ; max num secs to wait before SIGKILL (default 10)
user=root               ; setuid to this UNIX account to run the program
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/1
stderr_logfile_maxbytes=0


[program:celery_template_validations]
environment="NEW_RELIC_ENVIRONMENT=celery"
environment="NEW_RELIC_CONFIG_FILE=/etc/newrelic/newrelic_onyx.ini"
environment=
         C_FORCE_ROOT=true
command=bash -c "export NEW_RELIC_ENVIRONMENT=celery && export NEW_RELIC_CONFIG_FILE=/etc/newrelic/newrelic_onyx.ini && /usr/local/python3.8/bin/newrelic-admin run-program /usr/local/python3.8/bin/celery worker --time-limit=300 --concurrency=5 --app=celery_app --loglevel=INFO --maxtasksperchild=3000 -Q celery_template_validations"
directory=/usr/local/creditas/Onyx/onyx_proj/onyx_proj/
autostart=true
autorestart=true
startsecs=10                ; number of secs prog must stay running (def. 10)
startretries=3              ; max # of serial start failures (default 3)
exitcodes=1
stopsignal=TERM             ; signal used to kill process (default TERM)
stopwaitsecs=600             ; max num secs to wait before SIGKILL (default 10)
user=root               ; setuid to this UNIX account to run the program
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/1
stderr_logfile_maxbytes=0


[program:celery_heavy_data]
environment="NEW_RELIC_ENVIRONMENT=celery"
environment="NEW_RELIC_CONFIG_FILE=/etc/newrelic/newrelic_onyx.ini"
environment=
         C_FORCE_ROOT=true
command=bash -c "export NEW_RELIC_ENVIRONMENT=celery && export NEW_RELIC_CONFIG_FILE=/etc/newrelic/newrelic_onyx.ini && /usr/local/python3.8/bin/newrelic-admin run-program /usr/local/python3.8/bin/celery worker --time-limit=30000 --concurrency=1 --app=celery_app --loglevel=INFO --maxtasksperchild=3000 -n celery_heavy_data@%%h -Q celery_heavy_data"
directory=/usr/local/creditas/Onyx/onyx_proj/onyx_proj/
autostart=true
autorestart=true
startsecs=10                ; number of secs prog must stay running (def. 10)
startretries=3              ; max # of serial start failures (default 3)
exitcodes=1
stopsignal=TERM             ; signal used to kill process (default TERM)
stopwaitsecs=600             ; max num secs to wait before SIGKILL (default 10)
user=root               ; setuid to this UNIX account to run the program
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/1
stderr_logfile_maxbytes=0
